# contrib/bluetooth-timeout.service
[Unit]
Description=Bluetooth Timeout Daemon
After=network.target

[Service]
ExecStart=%h/.local/bin/bluetooth-timeout
WorkingDirectory=%h
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
